# miner-watcher.path
[Path]
Unit=miner-watcher.service
PathModified=/opt/backend

[Install]
WantedBy=multi-user.target